# chromeos-base/openjema-utils/files/init/jemaos-collect-log.conf
# Copyright 2023 Jema Innovations
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Collect jemaos log"
author          "John"

start on started system-services

oom score -100

exec /usr/share/cros/init/collect_jemaos_log.sh
